# backend/gunicorn.service
[Unit]
Description=SonicVision Gunicorn Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/SonicVision/backend
Environment="PATH=/root/SonicVision/backend/venv/bin"
EnvironmentFile=/root/SonicVision/backend/.env
ExecStart=/root/SonicVision/backend/venv/bin/gunicorn \
    --config /root/SonicVision/backend/gunicorn.conf.py \
    config.wsgi:application

# 日誌設置
StandardOutput=append:/var/log/sonicvision/gunicorn.log
StandardError=append:/var/log/sonicvision/gunicorn.error.log

# 重啟設置
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 