# backend/gunicorn.service
[Unit]
Description=SonicVision Gunicorn Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/sonicvision/backend
Environment="PATH=/var/www/sonicvision/backend/.venv/bin"
EnvironmentFile=/var/www/sonicvision/backend/.env.prod
ExecStart=/var/www/sonicvision/backend/.venv/bin/gunicorn \
    --config /var/www/sonicvision/backend/gunicorn.conf.py \
    config.wsgi:application

# 日誌設置
StandardOutput=append:/var/log/sonicvision/gunicorn.log
StandardError=append:/var/log/sonicvision/gunicorn.error.log

# 重啟設置
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 